# platform/broadcom/sonic-platform-modules-meituan/debian/platform-modules-meituan-ms25bcl-48y8c.postinst
depmod -a

systemctl enable platform-modules-meituan-ms25bcl-48y8c.service
systemctl start platform-modules-meituan-ms25bcl-48y8c.service

systemctl enable fanctrl.service
systemctl start fanctrl.service

/usr/local/bin/platform_api_mgnt.sh install
