# platform/broadcom/sonic-platform-modules-meituan/debian/platform-modules-meituan-ms100bcl-128c.postinst
depmod -a

systemctl enable platform-modules-meituan-ms100bcl-128c.service
systemctl start platform-modules-meituan-ms100bcl-128c.service

systemctl enable fanctrl.service
systemctl start fanctrl.service
